# systemd/fancontroller.socket
[Unit]
Description=fancontroller socket

[Socket]
ListenStream=/run/fancontroller/socket

[Install]
WantedBy=sockets.target